# util/ansible/roles/azuracast-cron/templates/supervisor.conf.j2
[program:cron]
command=/usr/sbin/cron -f
priority=600
numprocs=1
autostart=true
autorestart=unexpected
