# util/docker/final/service/php-fpm/run
#!/bin/bash

source /etc/php/.version

/usr/sbin/php-fpm${PHP_VERSION} -F --fpm-config /etc/php/${PHP_VERSION}/fpm/php-fpm.conf -c /etc/php/${PHP_VERSION}/fpm/
