# util/ansible/roles/redis/templates/supervisor.conf.j2
[program:redis]
command=/usr/bin/redis-server /etc/redis/redis.conf
user=redis
numprocs=1
autostart=true
autorestart=unexpected
