# util/docker/web/runit/php-worker/run
#!/bin/bash

sudo -E -u azuracast php /var/azuracast/www/bin/console queue:process 300 --worker-name=app_worker_0
