# util/docker/web/runit/php-worker/run
#!/bin/bash

sudo -E -u azuracast php /var/azuracast/www/util/cli.php queue:process 300