# util/docker/web/runit/beanstalkd/run
#!/bin/bash

echo 'Spinning up Beanstalkd process...'

sudo -E -u azuracast beanstalkd -l 127.0.0.1 -p 11300
