# util/docker/web/runit/beanstalkd/run
#!/bin/bash

echo 'Spinning up Beanstalkd process...'

exec sudo -E -u azuracast beanstalkd -p 11300 -z 262140
