# util/docker/web/runit/beanstalkd/run
#!/bin/bash

echo 'Spinning up Beanstalkd process...'

setuser azuracast beanstalkd -p 11300
