# util/docker/web/runit/sftpgo/run
#!/bin/bash

echo 'Spinning up SFTP process...'

cd /var/azuracast/sftpgo

sudo -E -u azuracast sftpgo --config-dir=/var/azuracast/sftpgo serve -l "" > /proc/1/fd/1 2> /proc/1/fd/2