# util/docker/web/runit/php-nowplaying/run
#!/bin/bash

svwaitup 5 /service/php-fpm || exit 1

exec gosu azuracast php /var/azuracast/www/bin/console azuracast:sync:nowplaying
