# util/ansible/roles/mariadb/templates/supervisor.conf.j2
[program:mariadb]
command=/usr/sbin/mariadbd
user=mysql
priority=100
numprocs=1
autostart=true
autorestart=unexpected
