# util/docker/web/service.full/acme.conf
[program:acme]
command=run_acme_sh
priority=200
numprocs=1
autostart=true
autorestart=true

stdout_logfile=/proc/1/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/proc/1/fd/2
stderr_logfile_maxbytes=0
